# rootfs/var/lib/dpkg/info/bb-customizations.list
/.
/etc
/etc/kernel
/etc/kernel/postinst.d
/etc/kernel/postinst.d/zz-uenv_txt
/etc/modprobe.d
/etc/modprobe.d/rtl8723bu-blacklist.conf
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/10-of-symlink.rules
/etc/udev/rules.d/60-omap-tty.rules
/etc/udev/rules.d/80-eeprom-noroot.rules
/etc/udev/rules.d/80-gpio-noroot.rules
/etc/udev/rules.d/80-i2c-noroot.rules
/etc/udev/rules.d/81-pwm-noroot.rules
/etc/udev/rules.d/82-gpio-config-pin.rules
/etc/udev/rules.d/83-eqep-noroot.rules
/etc/udev/rules.d/85-gpio-noroot.rules
/etc/udev/rules.d/86-remoteproc-noroot.rules
/etc/udev/rules.d/86-rpmsg-noroot.rules
/etc/udev/rules.d/87-iio-noroot.rules
/etc/udev/rules.d/88-leds-noroot.rules
/lib
/lib/firmware
/lib/firmware/am335x-pru0-fw.sleep
/lib/firmware/am335x-pru1-fw.sleep
/lib/firmware/am57xx-pru1_0-fw.sleep
/lib/firmware/am57xx-pru1_1-fw.sleep
/lib/firmware/am57xx-pru2_0-fw.sleep
/lib/firmware/am57xx-pru2_1-fw.sleep
/lib/systemd
/lib/systemd/system
/lib/systemd/system/generic-board-startup.service
/usr
/usr/share
/usr/share/doc
/usr/share/doc/bb-customizations
/usr/share/doc/bb-customizations/changelog.Debian.gz
/usr/share/doc/bb-customizations/copyright
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/ti_pru_firmware
