# rootfs/var/lib/dpkg/info/connman.conffiles
/etc/apport/crashdb.conf.d/connman-crashdb.conf
/etc/connman/main.conf
/etc/init.d/connman
