# rootfs/var/lib/dpkg/info/cpufrequtils.conffiles
/etc/init.d/cpufrequtils
/etc/init.d/loadcpufreq
