# rootfs/var/lib/dpkg/info/friendly-recovery.preinst
#!/bin/sh
set -e

# it is now a static service
systemctl disable friendly-recovery.service || :



exit 0
