# rootfs/var/lib/dpkg/info/lxcfs.list
/.
/etc
/etc/init.d
/etc/init.d/lxcfs
/etc/init
/etc/init/lxcfs.conf
/var
/var/lib
/var/lib/lxcfs
/usr
/usr/bin
/usr/bin/lxcfs
/usr/lib
/usr/lib/aarch64-linux-gnu
/usr/lib/aarch64-linux-gnu/lxcfs
/usr/lib/aarch64-linux-gnu/lxcfs/liblxcfs.so
/usr/lib/lxcfs
/usr/share
/usr/share/lxc
/usr/share/lxc/config
/usr/share/lxc/config/common.conf.d
/usr/share/lxc/config/common.conf.d/00-lxcfs.conf
/usr/share/lxcfs
/usr/share/lxcfs/lxc.mount.hook
/usr/share/lxcfs/lxc.reboot.hook
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/lxcfs.1.gz
/usr/share/doc
/usr/share/doc/lxcfs
/usr/share/doc/lxcfs/README.Debian
/usr/share/doc/lxcfs/copyright
/usr/share/doc/lxcfs/changelog.Debian.gz
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/lxcfs.service
/usr/lib/lxcfs/liblxcfs.so
