# rootfs/var/lib/dpkg/info/pollinate.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/pollinate
/usr/share/doc/pollinate/changelog.Debian.gz
/usr/share/doc/pollinate/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/pollinate.1.gz
/usr/bin
/usr/bin/pollinate
/etc
/etc/default
/etc/default/pollinate
/etc/pollinate
/etc/pollinate/entropy.ubuntu.com.pem
/etc/init
/etc/init/pollinate.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/pollinate.service
