# rootfs/var/lib/dpkg/info/ubuntu-advantage-tools.list
/.
/etc
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/20apt-esm-hook.conf
/etc/logrotate.d
/etc/logrotate.d/ubuntu-advantage-tools
/etc/ubuntu-advantage
/etc/ubuntu-advantage/help_data.yaml
/etc/ubuntu-advantage/uaclient.conf
/etc/update-manager
/etc/update-manager/release-upgrades.d
/etc/update-manager/release-upgrades.d/ubuntu-advantage-upgrades.cfg
/etc/update-motd.d
/etc/update-motd.d/88-esm-announce
/etc/update-motd.d/91-contract-ua-esm-status
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ua-messaging.service
/lib/systemd/system/ua-messaging.timer
/lib/systemd/system/ua-reboot-cmds.service
/usr
/usr/bin
/usr/bin/ubuntu-advantage
/usr/lib
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/uaclient
/usr/lib/python3/dist-packages/uaclient/__init__.py
/usr/lib/python3/dist-packages/uaclient/apt.py
/usr/lib/python3/dist-packages/uaclient/cli.py
/usr/lib/python3/dist-packages/uaclient/clouds
/usr/lib/python3/dist-packages/uaclient/clouds/__init__.py
/usr/lib/python3/dist-packages/uaclient/clouds/aws.py
/usr/lib/python3/dist-packages/uaclient/clouds/azure.py
/usr/lib/python3/dist-packages/uaclient/clouds/gcp.py
/usr/lib/python3/dist-packages/uaclient/clouds/identity.py
/usr/lib/python3/dist-packages/uaclient/config.py
/usr/lib/python3/dist-packages/uaclient/conftest.py
/usr/lib/python3/dist-packages/uaclient/contract.py
/usr/lib/python3/dist-packages/uaclient/defaults.py
/usr/lib/python3/dist-packages/uaclient/entitlements
/usr/lib/python3/dist-packages/uaclient/entitlements/__init__.py
/usr/lib/python3/dist-packages/uaclient/entitlements/base.py
/usr/lib/python3/dist-packages/uaclient/entitlements/cc.py
/usr/lib/python3/dist-packages/uaclient/entitlements/cis.py
/usr/lib/python3/dist-packages/uaclient/entitlements/esm.py
/usr/lib/python3/dist-packages/uaclient/entitlements/fips.py
/usr/lib/python3/dist-packages/uaclient/entitlements/livepatch.py
/usr/lib/python3/dist-packages/uaclient/entitlements/repo.py
/usr/lib/python3/dist-packages/uaclient/exceptions.py
/usr/lib/python3/dist-packages/uaclient/gpg.py
/usr/lib/python3/dist-packages/uaclient/pip.py
/usr/lib/python3/dist-packages/uaclient/security.py
/usr/lib/python3/dist-packages/uaclient/serviceclient.py
/usr/lib/python3/dist-packages/uaclient/snap.py
/usr/lib/python3/dist-packages/uaclient/status.py
/usr/lib/python3/dist-packages/uaclient/util.py
/usr/lib/python3/dist-packages/uaclient/version.py
/usr/lib/python3/dist-packages/ubuntu_advantage_tools-27.2.egg-info
/usr/lib/python3/dist-packages/ubuntu_advantage_tools-27.2.egg-info/PKG-INFO
/usr/lib/python3/dist-packages/ubuntu_advantage_tools-27.2.egg-info/dependency_links.txt
/usr/lib/python3/dist-packages/ubuntu_advantage_tools-27.2.egg-info/entry_points.txt
/usr/lib/python3/dist-packages/ubuntu_advantage_tools-27.2.egg-info/requires.txt
/usr/lib/python3/dist-packages/ubuntu_advantage_tools-27.2.egg-info/top_level.txt
/usr/lib/ubuntu-advantage
/usr/lib/ubuntu-advantage/apt-esm-hook
/usr/lib/ubuntu-advantage/apt-esm-json-hook
/usr/lib/ubuntu-advantage/patch_status_json.py
/usr/lib/ubuntu-advantage/reboot_cmds.py
/usr/lib/ubuntu-advantage/ua_update_messaging.py
/usr/lib/ubuntu-advantage/upgrade_lts_contract.py
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/ua
/usr/share/doc
/usr/share/doc/ubuntu-advantage-tools
/usr/share/doc/ubuntu-advantage-tools/changelog.gz
/usr/share/doc/ubuntu-advantage-tools/copyright
/usr/share/keyrings
/usr/share/keyrings/ubuntu-advantage-cis.gpg
/usr/share/keyrings/ubuntu-advantage-esm-apps.gpg
/usr/share/keyrings/ubuntu-advantage-esm-infra-trusty.gpg
/usr/share/keyrings/ubuntu-advantage-fips.gpg
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/ubuntu-advantage.1.gz
/var
/var/lib
/var/lib/ubuntu-advantage
/usr/bin/ua
/usr/share/man/man1/ua.1.gz
