# rootfs/usr/share/upstart/sessions/indicator-messages.conf
description "Indicator Messages Service"

start on indicator-services-start
stop on desktop-end or indicator-services-end

respawn
respawn limit 2 10

exec /usr/lib/aarch64-linux-gnu/indicator-messages/indicator-messages-service

