# rootfs/usr/share/upstart/xdg/autostart/indicator-printers.desktop
[Desktop Entry]
Type=Application
Name=Indicator Printers
Exec=/usr/lib/aarch64-linux-gnu/indicator-printers/indicator-printers-service
OnlyShowIn=Unity;
NoDisplay=true
StartupNotify=false
Terminal=false
Hidden=true
